# RemoteIPService/nas_root/etc/cron/cron.d/remoteip.sh
#! /bin/sh
wget http://www.g3ll3rt.de/remoteip/iptest.php
